# confmgr/files/confmgr.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2012 Xavier Carcelle <xavier.carcelle@gmail.com>

START=50
STOP=50

NAME=confmgr

start() {
	# import cache
	touch /tmp/confmgr.cache

	if [ -s /firmware/confmgr.cache ]; then
		confmgr fetch
	fi
}

stop() {
	# delete cache
	rm -f /tmp/confmgr.cache
}
